# tcp/roles/system_base/templates/os-updater.service.j2
[Unit]
Description=Safe unattended apt dist-upgrade (kernel-aware)
Documentation=man:apt-get(8)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/sbin/nitka-os-updater
Nice=10
TimeoutStartSec=2h
Environment=DEBIAN_FRONTEND=noninteractive
Environment=APT_LISTCHANGES_FRONTEND=none
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
