# tcp/roles/system_base/templates/os-updater.timer.j2
[Unit]
Description=Run nitka-os-updater nightly

[Timer]
OnCalendar=*-*-* 01:00
RandomizedDelaySec=13m
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
