# tcp/roles/ingress/templates/nitka-ingress-watchdog.timer.j2
[Unit]
Description=Run Nitka ingress Docker namespace watchdog fallback timer

[Timer]
OnBootSec=30s
OnUnitActiveSec={{ ingress_watchdog_timer_interval | default('5min') }}
AccuracySec=30s
Unit=nitka-ingress-watchdog.service

[Install]
WantedBy=timers.target
