# tcp/roles/ingress/templates/nitka-ingress-watchdog.service.j2
[Unit]
Description=Nitka ingress Docker namespace watchdog
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/nitka-ingress-watchdog --loop
Restart=always
RestartSec=5
TimeoutStartSec=30
Nice=5

[Install]
WantedBy=multi-user.target
